# templates/nfsen.service.erb
[Unit]
Description=NfSen Service
Requires=network.target remote-fs.target local-fs.target
After=network.target remote-fs.target local-fs.target

[Service]
Type=forking
PIDFile=<%= @_piddir %>/nfsend.pid
ExecStart=<%= @_bindir %>/nfsen start
ExecStop=<%= @_bindir %>/nfsen stop
Restart=on-abort

[Install]
WantedBy=multi-user.target
